# launchd/com.amplifier.slack-connector.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd service for the Amplifier Slack Connector.

  Setup:
    1. Copy this file to com.amplifier.slack-connector.plist
    2. Fill in your actual tokens in the EnvironmentVariables section
    3. Update paths to match your system
    4. Copy to ~/Library/LaunchAgents/com.amplifier.slack-connector.plist
    5. Load: launchctl load ~/Library/LaunchAgents/com.amplifier.slack-connector.plist
    6. Start: launchctl start com.amplifier.slack-connector

  Management:
    Status:  launchctl list com.amplifier.slack-connector
    Logs:    tail -f /tmp/slack-connector.log
    Stop:    launchctl stop com.amplifier.slack-connector
    Unload:  launchctl unload ~/Library/LaunchAgents/com.amplifier.slack-connector.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amplifier.slack-connector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/.local/share/uv/tools/amplifier/bin/slack-connector</string>
        <string>start</string>
        <string>--channel</string>
        <string>YOUR_CHANNEL_ID</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SLACK_BOT_TOKEN</key>
        <string>xoxb-your-bot-token-here</string>
        <key>SLACK_APP_TOKEN</key>
        <string>xapp-your-app-token-here</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>sk-ant-api03-your-api-key-here</string>
        <!-- Working directory for relative paths in bundle.md -->
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
        <!-- PATH: Use your full login shell PATH to ensure all tools are available -->
        <!-- Get it with: bash -l -c 'echo $PATH' -->
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>

    <!-- Keep alive: launchd will restart if the process exits -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart delay in seconds after crash -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Working directory (where bundle.md lives) -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/workspace/amplifier-module-connectors</string>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/tmp/slack-connector.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/slack-connector-error.log</string>
</dict>
</plist>
